# pi-image/files/byteorder-print.service
[Unit]
Description=ByteOrder Print Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=byteorder
WorkingDirectory=/opt/byteorder-print
EnvironmentFile=/etc/byteorder/print-service.env
ExecStart=/opt/byteorder-print/venv/bin/python -m app.main
Restart=on-failure
RestartSec=10
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
